# scripts/titan-gateway.service
[Unit]
Description=TITAN Gateway — Autonomous AI Agent Framework
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=dj
WorkingDirectory=/opt/TITAN
ExecStart=/usr/bin/node dist/cli/index.js gateway
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5
Environment=NODE_ENV=production
Environment=TITAN_GATEWAY_HOST=0.0.0.0
Environment=TITAN_PERSONAL=1
StandardOutput=append:/var/log/titan-gateway.log
StandardError=append:/var/log/titan-gateway.log

[Install]
WantedBy=multi-user.target
